# luci-theme-neobird/root/etc/uci-defaults/30_luci-theme-neobird
#!/bin/sh

exit 0 
# ↑↑↑ DO NOT EDIT THIS LINE
